# systemd/claude-tmux.service
[Unit]
Description=Persistent tmux session for managing Claude instances
Documentation=file://%h/src/dev/tmux-session/README.md
After=default.target

[Service]
Type=forking
# Start the statically-named session with the phone-friendly config.
# tmux daemonizes the server into this unit's cgroup, so the session (and
# every Claude instance spawned inside it) is owned by the user manager and
# survives logout (linger is enabled).
ExecStart=/usr/bin/tmux -f %h/src/dev/tmux-session/tmux.conf new-session -d -s claude
ExecStop=/usr/bin/tmux kill-session -t claude
RemainAfterExit=yes
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
